# infra/silo-control.service
[Unit]
Description=Silo Secure Terraform State Gateway - Control Plane
After=network.target vault.service

[Service]
Type=simple
User=root
WorkingDirectory=/etc/silo
ExecStart=/usr/local/bin/silo-server --component control-plane --config /etc/silo/silo.yaml
Restart=always
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
